# ansible/templates/networkshare.dns-wait.service.j2
[Unit]
Description=Wait for DNS resolution of SMB file server
After=systemd-resolved.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=120
ExecStart=/bin/sh -c 'h="{{ mount_what | default(ansible_local.pi_player.mount_what) | regex_replace("^//([^/]+)/.*$", "\\1") }}"; echo "$h" | grep -qE "^[0-9]+\.[0-9]+\.[0-9]+\.[0-9]+$" && exit 0; until getent hosts "$h" > /dev/null 2>&1; do sleep 2; done'
